# conf/supervisor.conf
[program:gunicorn]
user=l4s
directory=/srv/l4s/l4s
command=/srv/l4s/.virtualenvs/l4s/bin/gunicorn l4s.wsgi:application -c conf/gunicorn.conf
autostart=True
autorestart=True
redirect_stderr=True
